# slacko/roles/slacko/files/slacko-api.service
[Unit]
Description=Slack-O

[Service]
User=root
Group=root
ExecStart=/opt/slacko-api/slacko-api.sh

[Install]
WantedBy=multi-user.target